# inittab
:0:sysinit:/bin/mount -a
:0:respawn:/init
